# bin/org.openssh.sshd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>Label</key>
        <string>org.openssh.sshd</string>
        <key>OnDemand</key>
        <false/>
        <key>ProgramArguments</key>
        <array>
                <string>/opt/homebrew/sbin/sshd</string>
        </array>
        <key>ServiceDescription</key>
        <string>openssh server</string>
        <key>WorkingDirectory</key>
        <string>/opt/homebrew/sbin</string>
</dict>
</plist>
